# src/home/config/systemd/user/mbsync.timer
[Unit]
Description=mbsync query timer

[Timer]
OnBootSec=5s
OnUnitInactiveSec=30s

[Install]
WantedBy=default.target